# roles/postgres/files/pg_backup.timer
[Unit]
Description=PostgreSQL backup timer

[Timer]
OnCalendar=weekly
Persistent=true

[Wanted]
WantedBy=timers.target
