# roles/deploy/files/deploy.path
[Path]
PathChanged=/var/github
Unit=deploy.service

[Install]
WantedBy=multi-user.target
